# supervisord.conf
[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[supervisord]
pidfile=/run/supervisord.pid
nodaemon=true

[supervisorctl]
serverurl=unix:///run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:project_name]
numprocs=5
process_name=%(project_name)s_%(process_num)02d
command=python3 -m scrapy crawl %(process_num)
directory=/app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
